# etc/cron.d/noise
# m h dom mon dow user  command
30 19 * * * root /opt/sleepy-pi/noise.sh 12:00:00
